# etc/supervisord-serf.conf
[supervisord]
nodaemon=true

[program:serf-agent]
command=/docker-serf/bin/serf-agent.sh 
autorestart=true
