# etc/supervisord-haproxy.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/usr/sbin/haproxy -d -f /docker-serf/etc/haproxy.cfg
autorestart=true

